# nc-reporting.conf
[program:nc-reporting]
directory=/home/jnltech
command=python3 /home/jnltech/importQ.py
autostart=true
autorestart=true
